# uuu.conf
[program:uuu]
command=python3 /home/curator2/uuu/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/uuu/test.err.log
stdout_logfile=/home/curator2/uuu/test.out.log